# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/www/storage/logs/supervisord.log
pidfile=/var/run/supervisord.pid

[program:confirmation-consumer]
command=php /var/www/artisan kafka:consume-confirmation
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/confirmation-consumer.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

[program:inventory-consumer]
command=php /var/www/artisan kafka:consume-inventory
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/inventory-consumer.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

[program:analytics-consumer]
command=php /var/www/artisan kafka:consume-analytics
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/analytics-consumer.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

[program:dlq-monitor]
command=php /var/www/artisan kafka:consume-dlq
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/dlq-monitor.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
